# packages/network/iptables/system.d/iptables-defaults.service
[Unit]
Description=IPTables defaults
After=local-fs.target

ConditionPathExists=!/storage/.cache/services/iptables.conf
ConditionPathExists=!/storage/.cache/services/iptables.disabled

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'cp /usr/share/services/iptables.conf /storage/.cache/services/'
RemainAfterExit=yes
